# etc/build/dev/php7/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:php]
command = /usr/local/sbin/php-fpm
user = root
autostart = true

[program:credentials]
command = /root/credentials.sh
user = root
autostart = true
exitcodes=0,2

[program:sshd]
command=/usr/sbin/sshd -D
user = root
autostart = true
